# config_examples/proxy_py.supervisor.conf
[program:proxy_py]
command=/home/proxy_py/Env/proxy_py/bin/python3 /home/proxy_py/proxy_py/main.py 
directory=/home/proxy_py/proxy_py
user=proxy_py
autorestart=true
redirect_stderr=true
